# bootdebi/conf/defaults/xinetd.d/Xvnc
service Xvnc_1 
{ 
    type = UNLISTED
    disable = no 
    socket_type = stream 
    protocol = tcp
    wait = yes 
    user = root 
#    only_from = 10.0.0.0,127.0.0.1
    server = /usr/bin/Xvnc
    server_args = :1 -inetd -desktop localhost:1 -query localhost -auth /root/.Xauthority -depth 24 -rfbwait 30000 -rfbauth /etc/vnc/passwd -pn -SecurityTypes VncAuth
    port = 5901
}

